# site-cookbooks/infra-nfs/files/default/mkhome.cron
# This should run each 5 minutes to check if any home directory is missing.
*/5 * * * * root /opt/mkhome.sh
