# site-cookbooks/oe-infosec/files/default/aparser.cron
*/5 * * * *	cd /opt/oe-infosec/ && ./aparser.py
