# site-cookbooks/website2-dotcms/files/default/of.cron
*/5 * * * * root /opt/dotcms/lsofStats.sh >> /opt/dotcms/dotserver/tomcat/logs/of.log
